# com.jobbot.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS Launch Agent: auto-start FastAPI backend on login              -->
<!-- Install: cp com.jobbot.backend.plist ~/Library/LaunchAgents/         -->
<!--          launchctl load ~/Library/LaunchAgents/com.jobbot.backend.plist -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobbot.backend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>JOBBOT_DIR && source JOBBOT_DIR/.venv/bin/activate && cd JOBBOT_DIR && uvicorn backend.main:app --host 127.0.0.1 --port 8000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>JOBBOT_DIR</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>JOBBOT_DIR/data/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>JOBBOT_DIR/data/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
